# ubuntu/gridclicker.service
[Unit]
Description=GridClicker Daemon for Ubuntu
After=network.target

[Service]
Type=simple
User=%u
Environment=DISPLAY=:0
WorkingDirectory=/opt/gridclicker
ExecStart=/opt/gridclicker/venv/bin/python3 -u /opt/gridclicker/main.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
